# .kilo/worktrees/bottled-roundworm/config/supervisor/yalihan-video-render.conf
[program:yalihan-video-render]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/yalihan-panel/artisan queue:work --queue=video_render --tries=3 --timeout=600 --sleep=3 --max-jobs=1000 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/yalihan-panel/storage/logs/video-render-worker.log
stopwaitsecs=3600

