# .kilo/worktrees/righteous-tarascosaurus/launchd/com.yalihan.drive-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yalihan.drive-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>scripts/ops/sab-drive-sync.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/macbookpro/Projects/yalihan2026</string>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/macbookpro/Projects/yalihan2026/logs/drive-sync.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/macbookpro/Projects/yalihan2026/logs/drive-sync.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
